# pi/systemd/roverd.service
[Unit]
Description=Multi-Roomba rover control agent
After=network-online.target mediamtx.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/roverd -config /etc/roverd.yaml
Restart=on-failure
RestartSec=5
User=roverd
Group=roverd
AmbientCapabilities=CAP_SYS_TTY_CONFIG CAP_SYS_RAWIO

[Install]
WantedBy=multi-user.target
