# pi/systemd/video-publisher.service
[Unit]
Description=Rover Video Publisher (libcamera -> SRT)
After=network-online.target roverd.service
Wants=network-online.target

[Service]
Type=simple
User=roverd
Group=roverd
WorkingDirectory=/var/lib/roverd
EnvironmentFile=/var/lib/roverd/video.env
ExecStart=/usr/local/bin/video-publisher
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
