# roomcam-service/room-cam.service
[Unit]
Description=Room camera snapshot server (MJPEG webcam)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=DEVICE=/dev/video0
Environment=RESOLUTION=640x480
Environment=QUALITY=5
Environment=PORT=8088
Environment=WORKDIR=/run/roomcam
WorkingDirectory=/home/daniel/gits/MultiRoombaRover/roomcam-service
ExecStart=/usr/bin/env bash /home/daniel/gits/MultiRoombaRover/roomcam-service/room-cam-snapshot.sh
Restart=always
RestartSec=2
User=root
Group=root

[Install]
WantedBy=multi-user.target
